# 07-niord-backup/crontab
22 * * * * /hourly-backup.sh >> /var/log/cron.log 2>&1
25 06 * * * /daily-backup.sh >> /var/log/cron.log 2>&1

