# 07-niord-backup/crontab
25 * * * * /hourly-backup.sh >> /var/log/cron.log 2>&1
30 07 * * * /daily-backup.sh >> /var/log/cron.log 2>&1

